# cjdns/bin/utils/protomesh.service
[Unit]
Description=Manage Protomesh

[Service]
Type=oneshot
WorkingDirectory=/etc/proto-mesh/
ExecStart=/bin/bash /etc/proto-mesh/start.sh
ExecStop=/bin/bash /etc/proto-mesh/shutdown.sh
ExecReload=/bin/bash /etc/proto-mesh/shutdown.sh
ExecReload=/bin/bash /etc/proto-mesh/start.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
